# dist/ac-remote-celery.service
[Unit]
Description=ac-remote-celery worker
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/remote-ac-controller
ExecStart=/home/pi/remote-ac-controller/venv/bin/celery -A ac_ctl_web worker -B

ExecReload=/bin/kill -SIGUSR1 $MAINPID
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
