# dist/ac-remote-web.service
[Unit]
Description=ac-remote-web
After=network.target

[Service]
Type=simple
User=pi
Group=pi
Environment="PYTHONPATH=/home/app/.local/lib/python3.9/site-packages"
WorkingDirectory=/home/pi/remote-ac-controller
ExecStart=/home/pi/.virtualenvs/remote-ac-controller/bin/gunicorn ac_ctl_web.wsgi

ExecReload=/bin/kill -SIGUSR1 $MAINPID
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
